# sample-box/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:ttt-nginx]
command=/usr/local/bin/tuntuntun client -addr=http://host.docker.internal:8888 -transport=h2 -local-addr=:8080
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ttt-ssh]
command=/usr/local/bin/tuntuntun client -addr=http://host.docker.internal:8888 -transport=h2 -local-addr=:22
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
