# sample-box/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:iperf3]
command=iperf3 -s -p 5201
autostart=true
autorestart=true

[program:ttt]
command=/usr/local/bin/tuntuntun client -addr=http://host.docker.internal:8888 -transport=ws
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
